# home/.config/systemd/user/hypridle.service.d/override.conf
[Service]
Slice=background-graphical.slice
