# home/.config/systemd/user/hypridle.service.d/override.conf
[Service]
Slice=session-graphical.slice

# These values are experimental, your mileage may vary
CPUQuota=10%
MemoryHigh=10M
MemoryMax=50M
