# root/etc/systemd/system/fstrim.service.d/override.conf
[Unit]
ConditionACPower=true

[Service]
Nice=19
