# home/.config/systemd/user/wayland-wm@.service.d/override.conf
[Service]
# Reduce the likelihood of the compositor getting killed during OOM events
OOMScoreAdjust=-100
