# home/.config/systemd/user/hyprpaper.service.d/override.conf
[Service]
Slice=session-graphical.slice

# These values are experimental, your mileage may vary
CPUQuota=2%
MemoryHigh=50M
MemoryMax=100M
