# home/.config/systemd/user/hyprpaper.service.d/override.conf
[Service]
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "Hyprland" ""
Slice=background-graphical.slice
