# home/.config/systemd/user/hyprsunset.service.d/override.conf
[Service]
Slice=session-graphical.slice

# These values are experimental, your mileage may vary
CPUQuota=2%
MemoryHigh=3M
MemoryMax=15M
