# home/.config/systemd/user/wayland-wm@hyprland.desktop.service.d/override.conf
[Service]
# Guarantee this much minimum memory to Hyprland in case of OOM
# reclaim events, giving it an advantage against being killed.
MemoryLow=100M
