# root/etc/systemd/system/power-disconnected.target
[Unit]
Description=Power adapter is disconnected (charger NOT plugged)
Conflicts=power-connected.target

[Install]
WantedBy=multi-user.target
