# root/etc/systemd/system/power-connected.target
Description=Power adapter is connected (charger plugged)
Conflicts=power-disconnected.target

[Install]
WantedBy=multi-user.target
