# home/.config/systemd/user/eww-window@.service
[Unit]
Description=Open %I eww window
Documentation=https://elkowar.github.io/eww/
After=graphical-session.target
After=eww.service
PartOf=eww.service

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c "sleep 1 && eww ping"
ExecStart=/usr/bin/eww open %i
ExecStop=/usr/bin/eww close %i
RemainAfterExit=true
Restart=on-failure
RestartSec=3s
Slice=session-graphical.slice

[Install]
WantedBy=graphical-session.target
