# home/.config/systemd/user/eww.service
[Unit]
Description="Eww daemon"
PartOf=graphical-session.target

[Service]
ExecStart=/usr/local/bin/eww daemon --no-daemonize
Restart=on-failure

[Install]
WantedBy=wm-ready.target
