# home/.config/systemd/user/eww.service
[Unit]
Description=Eww daemon
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/eww daemon --no-daemonize
Restart=always
Slice=background-graphical.slice

[Install]
WantedBy=graphical-session.target
