# home/.config/systemd/user/nm-applet.service
[Unit]
Description=Network Manager Applet
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/nm-applet
Restart=on-failure

Slice=sessoin-graphical.slice

[Install]
WantedBy=graphical-session.target
