# home/.config/systemd/user/wm-ready.target
[Unit]
Description=Window Manager loaded and ready
BindsTo=graphical-session.target
